# root/etc/s6-overlay/s6-rc.d/service-rtorrent/run
#!/command/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

ARGS=""
if [[ -f "${CONFIG_DIR}/wireguard/forwarded_port" ]]; then
    port=$(cat "${CONFIG_DIR}/wireguard/forwarded_port")
    ARGS="-p ${port}-${port}"
fi

# shellcheck disable=SC2086
exec s6-setuidgid hotio "${APP_DIR}/rtorrent" -n -o "try_import=${CONFIG_DIR}/rtorrent.rc,network.scgi.open_local=/dev/shm/rtorrent.sock,system.daemon.set=true" ${ARGS}
